# server/grafana/grafana.service.template
[Unit]
Description=Grafana

[Service]
TimeoutStartSec=0
ExecStart=/home/monitoring/tools/grafana/start.sh

[Install]
WantedBy=multi-user.target