# server/grafana/grafana.service
[Unit]
Description=Grafana

[Service]
TimeoutStartSec=0
ExecStart=/home/monitoring/toolbox/grafana/start.sh

[Install]
WantedBy=multi-user.target
